# deployment/macos/com.phillippitts.speaktomack.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <!-- Service identification -->
    <key>Label</key>
    <string>com.phillippitts.speaktomack</string>

    <!-- Program and arguments -->
    <key>ProgramArguments</key>
    <array>
        <string>/usr/bin/java</string>
        <!-- JVM options -->
        <string>-Xms512m</string>
        <string>-Xmx2g</string>
        <string>-XX:+UseG1GC</string>
        <string>-XX:MaxGCPauseMillis=200</string>
        <string>-XX:+HeapDumpOnOutOfMemoryError</string>
        <string>-XX:HeapDumpPath=/var/log/speaktomack/heap-dump.hprof</string>
        <!-- Spring profile -->
        <string>-Dspring.profiles.active=production</string>
        <!-- Logging configuration -->
        <string>-Dlogging.config=/etc/speaktomack/log4j2-production.xml</string>
        <!-- JAR execution -->
        <string>-jar</string>
        <string>/opt/speaktomack/speaktomack.jar</string>
        <!-- Additional configuration -->
        <string>--spring.config.additional-location=/etc/speaktomack/application-production.properties</string>
    </array>

    <!-- Environment variables (absolute paths) -->
    <key>EnvironmentVariables</key>
    <dict>
        <key>VOSK_MODEL_PATH</key>
        <string>/opt/speaktomack/models/vosk-model-en-us-0.22</string>
        <key>WHISPER_MODEL_PATH</key>
        <string>/opt/speaktomack/models/ggml-base.en.bin</string>
        <key>WHISPER_BINARY_PATH</key>
        <string>/opt/speaktomack/bin/whisper</string>
        <key>SPEAKTOMACK_LOG_PATH</key>
        <string>/var/log/speaktomack</string>
    </dict>

    <!-- Working directory -->
    <key>WorkingDirectory</key>
    <string>/opt/speaktomack</string>

    <!-- Run at system start -->
    <key>RunAtLoad</key>
    <true/>

    <!-- Keep service alive (restart on crash) -->
    <key>KeepAlive</key>
    <true/>

    <!-- Logging -->
    <key>StandardOutPath</key>
    <string>/var/log/speaktomack/stdout.log</string>

    <key>StandardErrorPath</key>
    <string>/var/log/speaktomack/stderr.log</string>

    <!-- Process control -->
    <key>ThrottleInterval</key>
    <integer>10</integer>

    <!-- Resource limits -->
    <key>HardResourceLimits</key>
    <dict>
        <key>NumberOfFiles</key>
        <integer>65536</integer>
    </dict>

    <key>SoftResourceLimits</key>
    <dict>
        <key>NumberOfFiles</key>
        <integer>65536</integer>
    </dict>

    <!-- Security: don't run as root -->
    <key>UserName</key>
    <string>_speaktomack</string>
    <key>GroupName</key>
    <string>_speaktomack</string>

    <!-- Nice level (0 = normal priority) -->
    <key>Nice</key>
    <integer>0</integer>
</dict>
</plist>
